# scripts/machine-yearning.service
[Unit]
Description=Machine Yearning radio player
After=sound.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hal
Group=hal
WorkingDirectory=/home/hal/machine-yearning
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/hal/machine-yearning/.venv/bin/python -m player.radio --server
ExecStop=/bin/kill -TERM $MAINPID
KillMode=mixed
TimeoutStopSec=5
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
